# postgresql/rootfs/etc/services.d/mysqld/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
with-contenv
importas -i PGDATA PGDATA
s6-setuidgid postgres
postgres
