# mariadb/rootfs/etc/s6-overlay/s6-rc.d/set-root-user/up
/etc/s6-overlay/scripts/set-root-user.sh
