# base/rootfs/etc/s6-overlay/s6-rc.d/database-defaults/up
/etc/s6-overlay/scripts/database-defaults.sh
