# base/rootfs/etc/services.d/confd/run
#!/usr/bin/with-contenv bash
set -e 
confd-render-templates.sh --continuous
