# images/nginx/rootfs/etc/s6-overlay/s6-rc.d/environment-override-nginx/up
# Allow NGINX_ERROR_LOG_LEVEL to be overridden by DRUPAL_NGINX_ERROR_LOG_LEVEL, etc.
/usr/local/bin/confd-override-environment.sh --prefix NGINX
