# base/rootfs/etc/s6-overlay/s6-rc.d/confd-oneshot/up
/etc/s6-overlay/scripts/confd-oneshot.sh
