# handle/rootfs/etc/s6-overlay/s6-rc.d/handle/run
#!/usr/bin/env bash
set -e
exec s6-setuidgid handle /opt/handle/bin/hdl-server /var/handle
