# mariadb/rootfs/etc/s6-overlay/s6-rc.d/set-max-packet/up
/etc/s6-overlay/scripts/set-max-packet.sh
