# handle/rootfs/etc/s6-overlay/s6-rc.d/process-keys/up
/etc/s6-overlay/scripts/process-keys.sh
