# matomo/rootfs/etc/s6-overlay/s6-rc.d/matomo-setup/up
/etc/s6-overlay/scripts/matomo-setup.sh
