# images/drupal/rootfs/etc/s6-overlay/s6-rc.d/islandora-events-worker-index-fedora/run
#!/usr/bin/with-contenv bash
set -e
exec /etc/s6-overlay/scripts/messenger-worker.sh islandora_index_fedora "${DRUPAL_SM_WORKERS_INDEX_FEDORA_ENABLED:-true}"
