# images/handle/rootfs/etc/s6-overlay/s6-rc.d/set-db-driver/up
/etc/s6-overlay/scripts/set-db-driver.sh
