# code-server/rootfs/etc/s6-overlay/s6-rc.d/xdebug/run
#!/usr/bin/env bash
set -e

mkdir -p /run/xdebug

XDEBUG_FLAGS=$(</var/run/s6/container_environment/XDEBUG_FLAGS)

exec /usr/sbin/php-fpm83 \
	--pid /var/run/xdebug/php-fpm83.pid \
	--prefix /var/run/xdebug \
	--fpm-config /etc/php83/php-fpm.conf \
	-d "php_admin_value[zend_extension]=xdebug.so" \
	-d "zend_extension=xdebug.so" \
	"${XDEBUG_FLAGS[@]}"
