# karaf/rootfs/etc/services.d/karaf/run
#!/usr/bin/env bash
set -e
exec s6-setuidgid karaf /opt/karaf/bin/karaf server
