# karaf/rootfs/etc/services.d/karaf/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh :
with-contenv
s6-setuidgid karaf
/opt/karaf/bin/karaf server
