# nginx/rootfs/etc/s6-overlay/s6-rc.d/fpm-setup/up
# Create the folder to house the PID file and socket.
mkdir -p /run/php-fpm81
