# code-server/rootfs/etc/profile.d/path.sh
#!/usr/bin/env bash
export PATH=${PATH}:/var/lib/nginx/.composer/vendor/bin:/var/www/drupal/vendor/bin
