# images/tomcat/rootfs/etc/s6-overlay/s6-rc.d/environment-override-tomcat/up
# Allow TOMCAT_LOG_LEVEL to be overridden by FEDORA_TOMCAT_LOG_LEVEL, etc.
/usr/local/bin/confd-override-environment.sh --prefix TOMCAT
