# images/nginx/rootfs/etc/s6-overlay/s6-rc.d/fpm/run
#!/usr/bin/env bash
set -e
exec s6-notifyoncheck -d -t 30000 -n 1 /usr/sbin/php-fpm84 \
	--pid /run/php-fpm84/php-fpm84.pid \
	--prefix /run/php-fpm84 \
	--fpm-config /etc/php84/php-fpm.conf
