# mariadb/rootfs/etc/services.d/mysqld/run
#!/usr/bin/env bash
set -e
exec s6-setuidgid mysql mysqld --user mysql
