# mariadb/rootfs/etc/services.d/mysqld/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
s6-setuidgid mysql
mysqld --user mysql
