# activemq/rootfs/etc/services.d/activemq/run
#!/usr/bin/env bash
set -e
exec s6-setuidgid activemq /opt/activemq/bin/activemq console
