# base/rootfs/etc/s6-overlay/s6-rc.d/cacert-import/up
/etc/s6-overlay/scripts/cacert-import.sh
