# matomo/rootfs/etc/s6-overlay/s6-rc.d/set-subsite-defaults/up
/etc/s6-overlay/scripts/set-subsite-defaults.sh
