# images/postgresql/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/usr/bin/env bash
set -e
exec with-contenv \
	importas -i PGDATA PGDATA \
	s6-setuidgid \
	postgres postgres
