# handle/rootfs/etc/services.d/handle/finish
#!/usr/bin/execlineb -S1
# -*- mode: sh -*-
# vi: set ft=sh:
s6-svscanctl -t /var/run/s6/services
