# handle/rootfs/etc/services.d/handle/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
/opt/handle/bin/hdl-server /var/handle
