# handle/rootfs/etc/s6-overlay/s6-rc.d/handle-setup/up
/etc/s6-overlay/scripts/handle-setup.sh
