# images/riprap/rootfs/etc/s6-overlay/s6-rc.d/riprap-setup/up
/etc/s6-overlay/scripts/riprap-setup.sh
