# code-server/rootfs/etc/s6-overlay/s6-rc.d/development-environment/up
# Change uid of nginx user to match host.
/etc/s6-overlay/scripts/development-environment.sh
