# nginx/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/env bash
set -e
exec s6-notifyoncheck -t 30000 -n 1 /usr/sbin/nginx
