# nginx/rootfs/etc/s6-overlay/s6-rc.d/nginx-setup/up
# Create the folder to house the PID file.
mkdir -p /run/nginx
