# images/nginx/rootfs/etc/s6-overlay/s6-rc.d/environment-override-php/up
# Allow PHP_LOG_LEVEL to be overriden by DRUPAL_PHP_LOG_LEVEL, etc.
/usr/local/bin/confd-override-environment.sh --prefix PHP
