# nginx/rootfs/etc/s6-overlay/s6-rc.d/fpm/run
#!/usr/bin/env bash
set -e
exec s6-notifyoncheck -d -t 30000 -n 1 /usr/sbin/php-fpm82 \
	--pid /var/run/php-fpm82/php-fpm82.pid \
	--prefix /var/run/php-fpm82 \
	--fpm-config /etc/php82/php-fpm.conf
