# cantaloupe/rootfs/etc/s6-overlay/s6-rc.d/cantaloupe-defaults/up
/etc/s6-overlay/scripts/cantaloupe-defaults.sh
