# postgresql/rootfs/etc/s6-overlay/s6-rc.d/postgresql-setup/up
/etc/s6-overlay/scripts/postgresql-setup.sh
