# images/drupal/rootfs/etc/s6-overlay/s6-rc.d/islandora-events-worker-derivatives/finish
#!/command/execlineb -S2
/usr/local/share/s6/finish ${1} ${2}
