# images/drupal/rootfs/etc/s6-overlay/s6-rc.d/islandora-events-worker-derivatives/run
#!/usr/bin/with-contenv bash
set -e
exec /etc/s6-overlay/scripts/messenger-worker.sh islandora_derivatives "${DRUPAL_SM_WORKERS_DERIVATIVES_ENABLED:-true}"
