# code-server/rootfs/etc/services.d/xdebug/run
#!/usr/bin/env bash
set -e

XDEBUG_FLAGS=$(</var/run/s6/container_environment/XDEBUG_FLAGS)

exec /usr/sbin/php-fpm7 \
    --pid /var/run/xdebug/php-fpm7.pid \
    --prefix /var/run/xdebug \
    --fpm-config /etc/php7/php-fpm.conf \
    -d php_admin_value[zend_extension]=xdebug.so \
    -d zend_extension=xdebug.so \
    ${XDEBUG_FLAGS}