# nginx/rootfs/etc/services.d/fpm/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
/usr/sbin/php-fpm7
