# images/fcrepo6/rootfs/etc/s6-overlay/s6-rc.d/set-db-driver/up
# Some directories must exist prior to rendering templates.
/etc/s6-overlay/scripts/set-db-driver.sh
