# scyllaridae/rootfs/etc/s6-overlay/s6-rc.d/scyllaridae/run
#!/command/execlineb -P

with-contenv
exec s6-setuidgid scyllaridae /app/scyllaridae
