# nginx/rootfs/etc/services.d/nginx/run
#!/bin/sh
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
with-contenv /usr/sbin/nginx
