# tomcat/rootfs/etc/services.d/tomcat/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh :
with-contenv
s6-setuidgid tomcat
/opt/tomcat/bin/catalina.sh run
